# KEEDIO-AMBARI2/KEEDIO/1.4/services/GANGLIAUI/package/templates/gmond.service.j2
[Unit]
Description=Ganglia Monitoring Daemon
After=multi-user.target

[Service]
Type=forking
ExecStart=/usr/sbin/gmond -c /etc/ganglia/gmond.{{ name }}.conf

[Install]
WantedBy=multi-user.target

